# files/system/etc/systemd/user/plasma-kwin_wayland.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/kwin_wayland_wrapper