# config/files/usr/etc/system/upower.service.d/namespaces.conf
[Service]
# Namespaces
PrivateUsers=no
RestrictNamespaces=no