# files/system/usr/etc/profile.d/gnome-disable-jit.sh
export JavaScriptCoreUseJIT=0
export GJS_DISABLE_JIT=1
