# files/system/usr/etc/profile.d/gnome-disable-jit.sh
export JavaScriptCoreUseJIT=0 # See https://trac.webkit.org/wiki/EnvironmentVariables
export GJS_DISABLE_JIT=1 # See https://gitlab.gnome.org/GNOME/gjs/-/blob/master/doc/Environment.md
