# files/system/etc/systemd/user/org.gnome.Shell@wayland.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/gnome-shell --no-x11