# config/files/usr/etc/systemd/user/wireplumber.service.d/preload.conf
[Service]
Environment=LD_PRELOAD=/usr/lib64/libhardened_malloc-light.so
